# crontab/article.sh
#!/usr/bin/bash
cd /home/lucas/Graduation_Project/infoflow/ && python run.py crawl ifeng
cd /home/lucas/Graduation_Project/infoflow/ && python run.py crawl sougouWx
cd /home/lucas/Graduation_Project/infoflow/ && python run.py crawl qingboWx 
