# crontab/article.sh
#!/usr/bin/bash

while true
do
    cd /Users/gonglingxiao/Graduation_Project/spider && python run.py sougouweixin
	sleep 120
done
